# input/update-motd.d/20-tail
#!/bin/sh
#
# Last part of motd: /etc/motd.tail
#
# LOCATION      :/etc/update-motd.d/20-tail

[ -f /etc/motd.tail ] || exit 0

cat /etc/motd.tail
